# mcp_server/systemd/brain-maintenance.timer
[Unit]
Description=Run brain vault maintenance daily

[Timer]
# 04:00 local, a quiet hour; Persistent catches a missed run after downtime.
OnCalendar=*-*-* 04:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
